# launchd/com.signaldigest.app.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.signaldigest.app</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/python3</string>
        <string>/path/to/SimpleFeed/app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/SimpleFeed</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/path/to/SimpleFeed/logs/app.log</string>

    <key>StandardErrorPath</key>
    <string>/path/to/SimpleFeed/logs/app.err.log</string>
</dict>
</plist>
